# APIFrancescoJonusForNovitaDiamonds/docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm-entrypoint]
command=docker-php-entrypoint php-fpm